# deploy/hyperliquid-trader.service
[Unit]
Description=Hyperliquid Trading Scanner + Paper Trader
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/clawd/.openclaw/workspace/OpenClaw-Hyperliquid
EnvironmentFile=/home/clawd/.openclaw/workspace/hyperliquid-agent/.env
Environment=PKG_CONFIG_PATH=/home/linuxbrew/.linuxbrew/opt/openssl@3/lib/pkgconfig
ExecStart=/home/clawd/.openclaw/workspace/OpenClaw-Hyperliquid/target/release/hyperliquid-trader --paper-trade --interval 60 --starting-balance 110 --min-funding 20
Restart=always
RestartSec=10
StandardOutput=append:/home/clawd/.openclaw/workspace/OpenClaw-Hyperliquid/logs/trader.log
StandardError=append:/home/clawd/.openclaw/workspace/OpenClaw-Hyperliquid/logs/trader-error.log

[Install]
WantedBy=default.target
